# src/rc.local
#!/bin/bash
/etc/init.d/procps restart

exit 0